# repo/base/chakra-initscripts/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# set correct permissions for live user
/bin/chown -R live:users /home/live

